# packages/graphics/hyperion/system.d/hyperion.service
[Unit]
Description=Hyperion daemon
ConditionPathExists=/storage/.cache/services/hyperion.conf
After=plexht.service hyperion-defaults.service
Requires=hyperion-defaults.service

[Service]
ExecStartPre=/sbin/modprobe stk1160
ExecStart=/usr/bin/hyperiond /storage/.config/hyperion.config.json
ExecStopPost=/sbin/modprobe -r stk1160

[Install]
WantedBy=plexht.service
